# deploy/systemd/cloudbeds-rates-daily.service
[Unit]
Description=Cloudbeds Rates daily smoothing run
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=cloudbeds-rates
Group=cloudbeds-rates
WorkingDirectory=/opt/cloudbeds-rates
EnvironmentFile=/etc/cloudbeds-rates.env
ExecStart=/usr/bin/npm run daily:apply
TimeoutStartSec=6h
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
ReadWritePaths=/opt/cloudbeds-rates/data
